# scripts/pacman-update.timer
[Unit]
Description=Weekly pacman system update — runs after nightly backup

[Timer]
OnCalendar=Sun *-*-* 04:00:00
RandomizedDelaySec=30min
Persistent=true

[Install]
WantedBy=timers.target
